# packaging/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
	[ -f /var/lib/drlm/client.data.save ] && mv /var/lib/drlm/client.data.save /var/lib/drlm/client.data
	[ -f /var/lib/drlm/backup.data.save ] && mv /var/lib/drlm/backup.data.save /var/lib/drlm/backup.data
	[ -f /var/lib/drlm/network.data.save ] && mv /var/lib/drlm/network.data.save /var/lib/drlm/network.data
	[ -f /var/lib/drlm/drlm.sqlite.save ] && mv /var/lib/drlm/drlm.sqlite.save /var/lib/drlm/drlm.sqlite

  if [ $(ps -p 1 -o comm=) = "systemd" ]
  then
     systemctl enable tftpd-hpa.service
     systemctl enable nfs-kernel-server.service
     systemctl enable rpcbind.service
     systemctl enable isc-dhcp-server.service
     systemctl enable apache2.service
	   cp /usr/share/drlm/conf/systemd/drlm-stord.service /etc/systemd/system/
     systemctl daemon-reload
     systemctl enable drlm-stord.service
     systemctl start drlm-stord.service
  else
     update-rc.d tftpd-hpa defaults
     update-rc.d nfs-kernel-server defaults
     update-rc.d rpcbind defaults
     update-rc.d isc-dhcp-server defaults
     update-rc.d apache2 defaults
     cp /usr/sbin/drlm-stord /etc/init.d/
     update-rc.d drlm-stord defaults
     service drlm-stord start
  fi

  [ -f /etc/drlm/local.conf ] && echo "DHCP_SVC_NAME=\"isc-dhcp-server\"" >> /etc/drlm/local.conf
	;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
