# packaging/debian/postrm
#!/bin/sh
set -e

case "$1" in
  remove)
    # do something
    ;;

  install)
    # do something
    ;;

  upgrade)
    # do something
    ;;
  
  purge)
    # Stop, disable and remove drlm-stord.service
    if [ -f /etc/systemd/system/drlm-stord.service ]; then
      systemctl is-active --quiet drlm-stord.service && systemctl stop drlm-stord.service
      systemctl is-enabled --quiet drlm-stord.service && systemctl disable drlm-stord.service
      rm /etc/systemd/system/drlm-stord.service
    fi

    # Stop, disable and remove drlm-api.service
    if [ -f /etc/systemd/system/drlm-api.service ]; then
      systemctl is-active --quiet drlm-api.service && systemctl stop drlm-api.service
      systemctl is-enabled --quiet drlm-api.service && systemctl disable drlm-api.service
      rm /etc/systemd/system/drlm-api.service
    fi

    # Stop, disable and remove drlm-proxy.service
    if [ -f /etc/systemd/system/drlm-proxy.service ]; then
      systemctl is-active --quiet drlm-proxy.service && systemctl stop drlm-proxy.service
      systemctl is-enabled --quiet drlm-proxy.service && systemctl disable drlm-proxy.service
      rm /etc/systemd/system/drlm-proxy.service
    fi

    # Stop, disable and remove rsync service
    if [ -f /etc/systemd/system/drlm-rsyncd.service ]; then
      systemctl is-active --quiet drlm-rsyncd.service && systemctl stop drlm-rsyncd.service
      systemctl is-enabled --quiet drlm-rsyncd.service && systemctl disable drlm-rsyncd.service
      rm /etc/systemd/system/drlm-rsyncd.service
    fi

    # Stop, disable and remove tftpd service
    if [ -f /etc/systemd/system/drlm-tftpd.service ]; then
      systemctl is-active --quiet drlm-tftpd.service && systemctl stop drlm-tftpd.service
      systemctl is-enabled --quiet drlm-tftpd.service && systemctl disable drlm-tftpd.service
      rm /etc/systemd/system/drlm-tftpd.service 
    fi

    # Stop, disable and remove stunnel service
    if [ -f /etc/systemd/system/drlm-stunnel.service ]; then
      systemctl is-active --quiet drlm-stunnel.service && systemctl stop drlm-stunnel.service
      systemctl is-enabled --quiet drlm-stunnel.service && systemctl disable drlm-stunnel.service
      rm /etc/systemd/system/drlm-stunnel.service 
    fi

    # Stop, disable and remove gitd service
    if [ -f /etc/systemd/system/drlm-gitd.service ]; then
      systemctl is-active --quiet drlm-gitd.service && systemctl stop drlm-gitd.service
      systemctl is-enabled --quiet drlm-gitd.service && systemctl disable drlm-gitd.service
      rm /etc/systemd/system/drlm-gitd.service 
    fi

    systemctl daemon-reload

    # Remove all backups stored and archived
    rm -rf /var/lib/drlm 

    # Remove all client configurations
    rm -rf /etc/drlm

    # Remove all log files
    rm -rf /var/log/drlm

    ;;

  *)
    exit 0
    ;;

esac

exit 0
