# packaging/debian/postrm
#!/bin/sh
set -e

case "$1" in
    remove)
        # do something
	    ;;

    install)
        # do something
        ;;

    upgrade)
        # do something
	    ;;
    
    purge)
        rm -rf /var/lib/drlm 
        rm -rf /etc/drlm
        rm -rf /var/log/drlm

        if [ $(ps -p 1 -o comm=) = "systemd" ]
        then
            systemctl stop drlm-stord.service
            systemctl disable drlm-stord.service
            rm /etc/systemd/system/drlm-stord.service
            systemctl daemon-reload

            # systemctl stop tftpd-hpa.service
            # systemctl disable tftpd-hpa.service
            # systemctl stop nfs-kernel-server.service
            # systemctl disable nfs-kernel-server.service
            # systemctl stop rpcbind.service
            # systemctl disable rpcbind.service
            # systemctl stop isc-dhcp-server.service
            # systemctl disable isc-dhcp-server.service
        else
            service drlm-stord stop
            update-rc.d drlm-stord remove
            rm /etc/init.d/drlm-stord

            # update-rc.d tftpd-hpa stop
            # update-rc.d tftpd-hpa remove
            # update-rc.d nfs-kernel-server stop
            # update-rc.d nfs-kernel-server remove
            # update-rc.d rpcbind stop
            # update-rc.d rpcbind remove
            # update-rc.d isc-dhcp-server stop
            # update-rc.d isc-dhcp-server remove
        fi
	    ;;
    *)
        exit 0
        ;;
esac

exit 0
