# packaging/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove)
      if [ $(ps -p 1 -o comm=) = "systemd" ]
      then
        systemctl stop drlm-stord.service
        systemctl disable drlm-stord.service
        rm /etc/systemd/system/drlm-stord.service
        systemctl daemon-reload
      else
        service drlm-stord stop
        update-rc.d drlm-stord remove
        rm /etc/init.d/drlm-stord
      fi
      drlm_ver="$(awk 'BEGIN { FS="=" } /^VERSION=/ { print $$2}' /usr/sbin/drlm)"
      mv /var/lib/drlm/drlm.sqlite /var/lib/drlm/$drlm_ver-drlm.sqlite.save
      rm /etc/drlm/cert/drlm.*
      ;;

    upgrade)
      if [ $(ps -p 1 -o comm=) = "systemd" ]
      then
        systemctl stop drlm-stord.service
        systemctl disable drlm-stord.service
        rm /etc/systemd/system/drlm-stord.service
        systemctl daemon-reload
      else
        service drlm-stord stop
        update-rc.d drlm-stord remove
        rm /etc/init.d/drlm-stord
      fi
      drlm_ver="$(awk 'BEGIN { FS="=" } /^VERSION=/ { print $$2}' /usr/sbin/drlm)"
      mv /var/lib/drlm/drlm.sqlite /var/lib/drlm/$drlm_ver-drlm.sqlite.save
      ;;

    *)
      exit 0
      ;;
esac

exit 0
