# packaging/debian/prerm
#!/bin/sh
set -e

case "$1" in
  remove)
    # Stop, disable and remove drlm-stord.service
    systemctl is-active --quiet drlm-stord.service && systemctl stop drlm-stord.service
    systemctl is-enabled --quiet drlm-stord.service && systemctl disable drlm-stord.service
    [ -f /etc/systemd/system/drlm-stord.service ] && rm /etc/systemd/system/drlm-stord.service
    # Stop, disable and remove drlm-api.service
    systemctl is-active --quiet drlm-api.service && systemctl stop drlm-api.service
    systemctl is-enabled --quiet drlm-api.service && systemctl disable drlm-api.service
    [ -f /etc/systemd/system/drlm-api.service ] && rm /etc/systemd/system/drlm-api.service
    # Stop, disable and remove drlm-proxy.service
    systemctl is-active --quiet drlm-proxy.service && systemctl stop drlm-proxy.service
    systemctl is-enabled --quiet drlm-proxy.service && systemctl disable drlm-proxy.service
    [ -f /etc/systemd/system/drlm-proxy.service ] && rm /etc/systemd/system/drlm-proxy.service
    # Stop, disable and remove drlm-rsyncd.service
    systemctl is-active --quiet drlm-rsyncd.service && systemctl stop drlm-rsyncd.service
    systemctl is-enabled --quiet drlm-rsyncd.service && systemctl disable drlm-rsyncd.service
    [ -f /etc/systemd/system/drlm-rsyncd.service ] && rm /etc/systemd/system/drlm-rsyncd.service
    # Stop, disable and remove drlm-tftpd.service
    systemctl is-active --quiet drlm-tftpd.service && systemctl stop drlm-tftpd.service
    systemctl is-enabled --quiet drlm-tftpd.service && systemctl disable drlm-tftpd.service
    [ -f /etc/systemd/system/drlm-tftpd.service ] && rm /etc/systemd/system/drlm-tftpd.service
    # Stop, disable and remove drlm-stunnel.service
    systemctl is-active --quiet drlm-stunnel.service && systemctl stop drlm-stunnel.service
    systemctl is-enabled --quiet drlm-stunnel.service && systemctl disable drlm-stunnel.service
    [ -f /etc/systemd/system/drlm-stunnel.service ] && rm /etc/systemd/system/drlm-stunnel.service
    # Stop, disable and remove drlm-gitd.service
    systemctl is-active --quiet drlm-gitd.service && systemctl stop drlm-gitd.service
    systemctl is-enabled --quiet drlm-gitd.service && systemctl disable drlm-gitd.service
    [ -f /etc/systemd/system/drlm-gitd.service ] && rm /etc/systemd/system/drlm-gitd.service

    systemctl daemon-reload

    #Save the old database
    [ -f /usr/sbin/drlm ] && drlm_ver="$(awk 'BEGIN { FS="=" } /VERSION=/ { print $$2 }' /usr/sbin/drlm)"
    [ -f /var/lib/drlm/drlm.sqlite ] && mv /var/lib/drlm/drlm.sqlite /var/lib/drlm/$drlm_ver-drlm.sqlite.save
    rm /etc/drlm/cert/drlm.*
    ;;

  upgrade)
    # Stop, disable and remove drlm-stord.service
    systemctl is-active --quiet drlm-stord.service && systemctl stop drlm-stord.service
    systemctl is-enabled --quiet drlm-stord.service && systemctl disable drlm-stord.service
    [ -f /etc/systemd/system/drlm-stord.service ] && rm /etc/systemd/system/drlm-stord.service
    # Stop, disable and remove drlm-api.service
    systemctl is-active --quiet drlm-api.service && systemctl stop drlm-api.service
    systemctl is-enabled --quiet drlm-api.service && systemctl disable drlm-api.service
    [ -f /etc/systemd/system/drlm-api.service ] && rm /etc/systemd/system/drlm-api.service
    # Stop, disable and remove drlm-proxy.service
    systemctl is-active --quiet drlm-proxy.service && systemctl stop drlm-proxy.service
    systemctl is-enabled --quiet drlm-proxy.service && systemctl disable drlm-proxy.service
    [ -f /etc/systemd/system/drlm-proxy.service ] && rm /etc/systemd/system/drlm-proxy.service
    # Stop, disable and remove drlm-rsyncd.service
    systemctl is-active --quiet drlm-rsyncd.service && systemctl stop drlm-rsyncd.service
    systemctl is-enabled --quiet drlm-rsyncd.service && systemctl disable drlm-rsyncd.service
    [ -f /etc/systemd/system/drlm-rsyncd.service ] && rm /etc/systemd/system/drlm-rsyncd.service
    # Stop, disable and remove drlm-tftpd.service
    systemctl is-active --quiet drlm-tftpd.service && systemctl stop drlm-tftpd.service
    systemctl is-enabled --quiet drlm-tftpd.service && systemctl disable drlm-tftpd.service
    [ -f /etc/systemd/system/drlm-tftpd.service ] && rm /etc/systemd/system/drlm-tftpd.service
    # Stop, disable and remove drlm-stunnel.service
    systemctl is-active --quiet drlm-stunnel.service && systemctl stop drlm-stunnel.service
    systemctl is-enabled --quiet drlm-stunnel.service && systemctl disable drlm-stunnel.service
    [ -f /etc/systemd/system/drlm-stunnel.service ] && rm /etc/systemd/system/drlm-stunnel.service
    # Stop, disable and remove drlm-gitd.service
    systemctl is-active --quiet drlm-gitd.service && systemctl stop drlm-gitd.service
    systemctl is-enabled --quiet drlm-gitd.service && systemctl disable drlm-gitd.service
    [ -f /etc/systemd/system/drlm-gitd.service ] && rm /etc/systemd/system/drlm-gitd.service

    systemctl daemon-reload

    #Save the old database
    [ -f /usr/sbin/drlm ] && drlm_ver="$(awk 'BEGIN { FS="=" } /VERSION=/ { print $$2 }' /usr/sbin/drlm)"
    [ -f /var/lib/drlm/drlm.sqlite ] && mv /var/lib/drlm/drlm.sqlite /var/lib/drlm/$drlm_ver-drlm.sqlite.save
    ;;

  purge)
    # Unconfigure nbd
    /usr/share/drlm/conf/nbd/config-nbd.sh remove
    ;;

  *)
    exit 0
    ;;

esac

exit 0
