# etc/cron.d/drlm
* * * * * root /usr/sbin/drlm sched --run
